# debian/dfmplugin-disk-encrypt.postinst
#!/bin/sh

echo "Enable dde-filemanager-daemon.service"
systemctl enable dde-filemanager-daemon.service || true
